# root/etc/services.d/announcedd/run
#!/usr/bin/with-contenv bash

umask 022
cd /usr/share/announcedd
exec s6-setuidgid abc /usr/bin/python3 /usr/share/announcedd/bot.py
